# wimazza.timer
[Unit]
Description=WiMazza timer login

[Timer]
OnUnitActiveSec=1min
OnActiveSec=5s

[Install]
WantedBy=timers.target
