# wimazza.service
[Unit]
Description=WiMazza oneshot login

[Service]
Type=oneshot
ExecStart=/bin/bash /home/baronefr/.wimazzer/wimazzer.sh
